# deploy/yamparser.service
[Unit]
Description=Yamparser - парсер отзывов из Google Sheets
After=network.target

[Service]
Type=simple
User=danil
WorkingDirectory=/home/danil/yamparser
ExecStart=/usr/bin/xvfb-run /home/danil/yamparser/venv/bin/python integrated_parser.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yamparser

[Install]
WantedBy=multi-user.target
